# rc.local
#!/bin/sh -e

sudo -H -u vagrant dtach -n /tmp/rtsess rtorrent
